# services/wolsimpleserver.target
[Unit]
Description=WOL Simple Server
Requires=wolsimpleserver_main.service wolsimpleserver_proxy.service

[Install]
WantedBy=multi-user.target